# units/resize/resize.service
[Unit]
Description=Video resize

[Service]
Type=oneshot
ExecStart=/usr/local/bin/resize.sh

[Install]
WantedBy=multi-user.target
